# src/services/jobs/postgresql_node_ng/config/postgresql_logrotate.cron

*/5 * * * * test -x /usr/sbin/logrotate && /usr/sbin/logrotate /var/vcap/jobs/postgresql_node_ng/config/logrotate.conf >>/var/vcap/sys/log/postgresql_node/postgresql_logrotate_cron.log 2>&1
