# src/services/jobs/postgresql_node_ng/config/postgresql_backup.cron


0 7 * * * /var/vcap/packages/ruby/bin/ruby /var/vcap/packages/postgresql_node_ng/services/ng/postgresql/bin/postgresql_backup -c /var/vcap/jobs/postgresql_node_ng/config/postgresql_backup.yml >>/var/vcap/sys/log/postgresql_node/postgresql_backup_cron.log 2>&1
